# linux/Eos-rankmirrors/eos.service
[Unit]
Description=Update Endeavour mirrorlist using reflector

[Service]
Type=oneshot
ExecStart=/usr/bin/eos-rankmirrors
